# pi/systemd/aa-clock-monitor.service
[Unit]
Description=Android Auto Clock Monitor
After=openauto.service

[Service]
Type=simple
ExecStart=/home/pi/aa_clock_monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
